# launch-agents/Library/LaunchAgents/com.sulami.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sulami.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/sulami/.restic/backup.rb</string>
    <string>backup-and-cleanup</string>
  </array>
  <key>ExitTimeOut</key>
  <integer>120</integer>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>19</integer>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
